# app/utils/system/jobs/crontab
*/5 * * * *  cd /app && /usr/local/bin/python3 /app/get_btc_data.py  > /proc/1/fd/1 2>/proc/1/fd/2
*/15 * * * *  cd /app && /usr/local/bin/python3 /app/meme_index.py  > /proc/1/fd/1 2>/proc/1/fd/2
